# assets/systemd/hyprdiscover-check.timer
[Unit]
Description=HyprDiscover periodic update check

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
Persistent=true

[Install]
WantedBy=default.target
