# assets/systemd/hyprdiscover-check.service
[Unit]
Description=HyprDiscover update check
After=network-online.target

[Service]
Type=oneshot
ExecStart=hyprdiscover --check
StandardOutput=journal
StandardError=journal
